# system_files/etc/profile.d/nimbus.sh
alias neofetch='/usr/bin/fastfetch'
alias njust='/usr/bin/just -f /usr/share/nimbus/justfile'